# pi_server/scripts/rover.service
[Unit]
Description=Rover Control Server (WiFi/Web UI)
Documentation=https://github.com/Surrealabs/Ogun
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rover_server /etc/rover/rover.conf
Restart=on-failure
RestartSec=5
User=root
Group=root
WorkingDirectory=/opt/rover

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rover

# Resource limits
LimitRTPRIO=5
LimitNICE=-5

[Install]
WantedBy=multi-user.target
