# scripts/rover-auto-update.service
[Unit]
Description=Ogun Rover Auto Updater
Documentation=https://github.com/Surrealabs/Ogun
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rover-auto-update.sh
User=root
Group=root

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=rover-auto-update
